# deploy/supervisor-segob-queue.conf
; Supervisor — worker de cola para exportaciones Excel / Word (módulos temporales)
; Instalación (Linux):
;   sudo cp deploy/supervisor-segob-queue.conf /etc/supervisor/conf.d/segob-queue.conf
;   sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start segob-queue:*
;
; Ajusta user y directory a tu servidor.

[program:segob-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/segob/artisan queue:work database --sleep=3 --tries=2 --timeout=3600 --max-jobs=50
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/segob/storage/logs/queue-worker.log
stopwaitsecs=3600
